# resolver/pkg/scripts/postinst
#!/bin/sh
set -e

if ! getent passwd pzresolver >/dev/null 2>&1; then
    adduser --system --group --no-create-home --home /nonexistent pzresolver
fi

install -d -o root       -g root       -m 0755 /etc/promtuz
install -d -o root       -g root       -m 0711 /etc/promtuz/keys
install -d -o root       -g root       -m 0755 /etc/promtuz/certs
install -d -o pzresolver -g pzresolver -m 0700 /etc/promtuz/keys/resolver

systemctl daemon-reload >/dev/null 2>&1 || true

# On upgrade, pick up the new binary immediately. try-restart only acts if the
# service is already running, so it's a no-op on a fresh (not-yet-enabled)
# install — the operator still runs `enable --now` per the steps below.
systemctl try-restart pzresolver >/dev/null 2>&1 || true

cat <<'EOF'
pzresolver installed. Next:
  1. systemctl enable --now pzresolver     # generates the key, writes the CSR, waits
  2. certgen sign /etc/promtuz/keys/resolver/resolver.csr  (on your CA box)
  3. copy the signed cert to /etc/promtuz/certs/resolver.crt  # starts automatically
EOF

exit 0